# ingestor-promo/crontab
*/10 * * * * /usr/local/bin/python3 /app/fetch_promo.py >> /var/log/cron.log 2>&1
