# ingestor-sql/crontab
*/10 * * * * /usr/local/bin/python3 /app/export_sql.py >> /var/log/cron.log 2>&1
